# var/etc/systemd/system/viper-gunicorn.socket
[Unit]
Description=Viper Gunicorn Socket

[Socket]
ListenStream=/srv/src/dj_viper/var/run/gunicorn.sock

[Install]
WantedBy=sockets.target
